# docker-alluxio-base/supervisord-slave.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:sshd]
command=/usr/sbin/sshd

[program:alluxio-slave]
command=bash -c '/alluxio/bin/alluxio-start.sh worker SudoMount'
