# docker-alluxio-base/supervisord-local.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:sshd]
command=/usr/sbin/sshd

[program:alluxio-local]
command=bash -c '/alluxio/bin/alluxio-start.sh local'
